# etc/profile.d/home-cargo-bin.sh
export PATH="$HOME/.local/cargo/bin:$PATH"