# etc/profile.d/program.sh
export EDITOR=nvim
export VISUAL=nvim
export PAGER=less
export DIFFPROG=delta
export DOCKER_BUILDKIT=1
export GO111MODULE=on
export GOPROXY=https://goproxy.io,direct
